# ngrok/rootfs/etc/services.d/ngrok/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting ngrok..."

exec ngrok start --config /ngrok-config/ngrok.yml --all
